# systemd/spotify-kiosk-watchdog.service
[Unit]
Description=Spotify kiosk watchdog (probe DevTools + Spotify Connect device)

[Service]
Type=oneshot
ExecStart=%h/bin/spotify-auto watchdog
