# systemd/spotify-kiosk-nightly-restart.service
[Unit]
Description=Nightly preventive restart of spotify-kiosk.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart spotify-kiosk.service
